# scripts/systemd/check-single-instance.timer
[Unit]
Description=Ежедневный аудит дублей инстансов ботов

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
